# ops/systemd/claudeusageproxy.service
[Unit]
Description=Claude usage proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/mudro/app
EnvironmentFile=-/etc/openclaw/runtime/claudeusageproxy.env
Environment=HOME=/home/openclaw
ExecStart=/opt/mudro/app/scripts/openclaw/claudeusageproxy_systemd.sh
Restart=on-failure
RestartSec=5
UMask=0027
NoNewPrivileges=yes
ReadWritePaths=/var/lib/openclaw /var/log/openclaw /home/openclaw

[Install]
WantedBy=multi-user.target
