# ops/systemd/mudro-agent-planner.service
[Unit]
Description=MUDRO agent planner (oneshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=mudro
Group=mudro
SupplementaryGroups=docker
WorkingDirectory=/opt/mudro/app
EnvironmentFile=-/etc/mudro/runtime/mudro-agent.env
Environment=HOME=/var/lib/mudro
Environment=XDG_CACHE_HOME=/var/lib/mudro/.cache
Environment=GOCACHE=/var/lib/mudro/go-build
ExecStart=/opt/mudro/bin/mudro-agent --mode planner-once
UMask=0027
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
ReadWritePaths=/var/lib/mudro /var/log/mudro /opt/mudro/app/.codex /var/run/docker.sock

[Install]
WantedBy=multi-user.target
