# ops/systemd/openclaw.service
[Unit]
Description=OpenClaw gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/mudro/app
EnvironmentFile=-/etc/openclaw/runtime/openclaw.env
Environment=HOME=/home/openclaw
ExecStart=/opt/mudro/app/scripts/openclaw/openclaw_gateway_systemd.sh
Restart=on-failure
RestartSec=5
UMask=0027
NoNewPrivileges=yes
ReadWritePaths=/var/lib/openclaw /var/log/openclaw /home/openclaw

[Install]
WantedBy=multi-user.target
