# ops/systemd/skaro.service
[Unit]
Description=Skaro dashboard
After=network-online.target claudeusageproxy.service
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/mudro/app
EnvironmentFile=-/etc/openclaw/runtime/skaro.env
Environment=HOME=/home/openclaw
ExecStart=/opt/mudro/app/scripts/skaro/skaro_ui_linux.sh
Restart=on-failure
RestartSec=5
UMask=0027
NoNewPrivileges=yes
ReadWritePaths=/var/lib/openclaw /var/log/openclaw /home/openclaw

[Install]
WantedBy=multi-user.target
